# deploy/supervisord.conf
[supervisord]

[supervisorctl]
serverurl = unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[inet_http_server]
port = 127.0.0.1:9001

[unix_http_server]
file = supervisor.sock

[program:create_dir]
command = mkdir -m 777 -p /var/run/postgresql
autorestart=false

[program:nginx]
command = nginx -c %(ENV_FLESH_DEPLOY_DIR)s/nginx.conf

[program:postgres]
command = su postgres -c "%(ENV_PG)s/postgres -D /etc/postgresql/%(ENV_PG_VERSION)s/flesh -p 5454"

[program:pgpool]
command = pgpool -n

[program:revel]
directory = %(ENV_FLESHLOCATION)s
command = revel run flesh prod

[program:deploy_listener]
# listen on port 9899 for any connection, then run the deploy script and quit
command = /bin/nc.traditional -l -p 9899 -c "%(ENV_FLESH_DEPLOY_DIR)s/deploy.sh" localhost